# dist/celerybeat.conf
; ================================
;  celery beat supervisor example
; ================================

[program:celerybeat_<add_app_name>]
environment=PYTHONPATH="<path_to_nemesis>:<path_to_hippo/hippocates>"
command=<celery_from_virtualenv> -A hippocrates.celery_worker beat --loglevel=INFO
directory=<path_to_hippo>

user=tmis
numprocs=1
stdout_logfile=/var/log/celery/%(program_name)s_beat.log
stderr_logfile=/var/log/celery/%(program_name)s_beat.log
autostart=true
autorestart=true
startsecs=10

; if rabbitmq (or other broker) is supervised, set its priority higher
; so it starts first
priority=1001
